# deployment/systemd/worldcup-update-pick-results.timer
[Unit]
Description=Run Daily Picks Result Updater at 11pm

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
